# deploy/llm/llm-warmup.timer.example
[Unit]
Description=SecWay LLM alle 2 Minuten warmhalten
[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
[Install]
WantedBy=timers.target
