# deploy/llm/llm-warmup.service.example
[Unit]
Description=SecWay LLM warmhalten
After=llm-classify.service
[Service]
Type=oneshot
User=llm
ExecStart=/opt/llm/warmup.sh
